# roles/verify_system/templates/patchmgmt.timer.j2
[Unit]
Description=Run Patch Management ones a week

[Timer]
OnCalendar=Mon *-*-* 00:30:00
Persistent=true

[Install]
WantedBy=timers.target