# roles/verify_system/templates/system_update_checker.timer.j2
[Unit]
Description=Start system update checker

[Timer]
OnBootSec=5min
Persistent=true

[Install]
WantedBy=timers.target